# app/crontab
*/1 * * * * python3 /app/manage.py course-holding-maker
